# etc/cron.d/watchdog-reset.cron
* * * * * root /usr/local/sbin/watchdog-reset
